# deploy/linux/doc-review.service
[Unit]
Description=Doc Review API (FastAPI + Uvicorn)
After=network.target postgresql.service
Wants=postgresql.service
# 主服务启动时一并拉起 1 个 OCR worker(独立进程,OCR 崩溃不连累 API)。
# 需要更多 worker:手动 systemctl enable --now doc-review-worker@2 ...
Wants=doc-review-worker@1.service

[Service]
Type=simple
User=docreview
Group=docreview
WorkingDirectory=/opt/doc-review/backend
Environment=PATH=/opt/doc-review/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=LANG=C.UTF-8
Environment=LC_ALL=C.UTF-8
Environment=PYTHONUTF8=1
Environment=OMP_NUM_THREADS=2
Environment=MKL_NUM_THREADS=2
Environment=ARCHIVE_DETECT_QUEUE_MAX=200
# 数据库:取 EnvironmentFile 里的 DATABASE_URL,优先级高于 config.json
EnvironmentFile=-/opt/doc-review/deploy/linux/app.env

# 主进程(方案二 2b)只接 HTTP / 写 DB / 跑 finalize 协程 / watchdog 回收死 worker 任务
# OCR Worker 是独立 systemd 服务: doc-review-worker@{1,2,3}.service
# --workers=1 因为 uvicorn 主进程内有内存态 _batch_status 缓存(前端轮询 fast-path)
ExecStart=/opt/doc-review/.venv/bin/python -m uvicorn main:app \
    --host 127.0.0.1 --port 8000 --workers 1 \
    --proxy-headers --forwarded-allow-ips=127.0.0.1

Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# 资源限制(主进程很轻,2GB 足够)
MemoryMax=2G
MemoryHigh=1.5G
LimitNOFILE=65535

# 沙箱:
# - NoNewPrivileges 防提权
# - PrivateTmp 给本进程独立 /tmp(LibreOffice 多次启动不互锁)
# - ProtectSystem=full 让 /usr/etc 只读
# - ProtectHome=no:PaddleOCR 要写 ~/.paddlex/(模型缓存),LibreOffice 要写 ~/.config/libreoffice/(profile)
#   不能用 read-only,会让 OCR 首次初始化报错
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=no
ReadWritePaths=/opt/doc-review/output /opt/doc-review/temp

# 日志走 journald,自带 rotation + ratelimit,不再写本地文件避免无限增长
StandardOutput=journal
StandardError=journal
SyslogIdentifier=doc-review

[Install]
WantedBy=multi-user.target
